# configs/subs-api-v6.conf
[program:subs-api-v6]
command=/data/node/go/bin/subs-api-v6
autostart=true
user=node
redirect_stderr=true
stdout_logfile=/data/log/node/subs-api-v6.log
stdout_logfile_maxbytes=10MB
environment=HOME="/data/node"
directory=/data/node/go
